# .ProjectCore/Configuration/Docker/PHPApacheBackEnd/System/etc/cron.d/cron_Script_ImportFile_SCP_SwingBarrierGate_sh
*/5 * * * * root /zhtConf/Script/BashScript/Script.ImportFile.SCP.SwingBarrierGate.sh >> /var/log/cron.log 2>&1
# Don’t remove the empty line at the end of this file. It is required to run the cron job
