# .ProjectCore/Configuration/Docker/PHPApacheBackEnd/System/etc/cron.d/cron_Script_Laravel_ArtisanScheduleRun_sh
* * * * * root /zhtConf/Script/BashScript/Script.Laravel.ArtisanScheduleRun.sh >> /var/log/cron.log 2>&1
# /var/www/html/WebBackEnd/storage/logs/cron.log
# Don’t remove the empty line at the end of this file. It is required to run the cron job
